# scripts/market-sniper.service
[Unit]
Description=Market Sniper Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/US-trading-Crypto/backend
ExecStart=/bin/bash /opt/US-trading-Crypto/scripts/run_sniper.sh
Restart=always
RestartSec=10

# Logging - all output goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=market-sniper

# Keep process alive
KillMode=process
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
