# scripts/market-sniper.service
[Unit]
Description=Market Sniper Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/US-trading-Crypto/backend
ExecStart=/bin/bash /opt/US-trading-Crypto/scripts/run_sniper.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

