# scripts/market-sniper.service
[Unit]
Description=Market Sniper Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/US-trading-Crypto
ExecStart=/usr/bin/python3 /opt/US-trading-Crypto/backend/market_sniper.py
Restart=always
RestartSec=10
Environment="PATH=/usr/bin:/usr/local/bin"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
